# scripts/udev/30-docking.rules
# rules for docking / undocking
ENV{EVENT}=="undock", KERNEL=="dock.0", RUN+="dock"
ENV{EVENT}=="dock", KERNEL=="dock.0", RUN+="dock"
